# scripts/systemd/ssh-key-audit.path
[Unit]
Description=ec2-user authorized_keys degisikligini izle -> sertlestirme denetimi

[Path]
PathModified=/home/ec2-user/.ssh/authorized_keys
Unit=ssh-key-audit.service

[Install]
WantedBy=multi-user.target
